# root/etc/s6/nginx/finish
#!/bin/bash
echo "nginx finish script runned"
