# root/etc/s6/nginx/run
#!/bin/bash
cd /etc/nginx
. /create-nginx-sites.sh
nginx
