# root/etc/s6/puppetdb/finish
#!/bin/bash
echo "php5-fpm finish script runned"
